# docker/root/etc/services.d/brultech-serial2mqtt/run
#!/usr/bin/with-contenv bash

cd /src/brultech_serial2mqtt || exit
exec \
    s6-setuidgid abc python3.11 -m brultech_serial2mqtt
