# docker/root/etc/services.d/brultech-serial2mqtt/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc python3 -m brultech_serial2mqtt
